# deploy/forex-bot.service
[Unit]
Description=Forex trend-pullback bot (EURUSD/GBPUSD)
Documentation=file:///opt/forex-bot/README.md
After=network-online.target mt5-bridge.service
Wants=network-online.target mt5-bridge.service
# Soft dependency: bot retries MT5 on its own; do not fail boot if bridge is still waiting on IPC.

[Service]
Type=simple
User=forexbot
Group=forexbot
WorkingDirectory=/opt/forex-bot
EnvironmentFile=/opt/forex-bot/.env
# Drop --paper once paper mode is validated.
# --db must point inside a ReadWritePaths dir: ProtectSystem=strict makes the
# rest of /opt/forex-bot read-only, so the default ./forex_bot.db would fail.
ExecStart=/opt/forex-bot/.venv/bin/python /opt/forex-bot/run_live.py --paper --db data/forex_bot.db
Restart=on-failure
RestartSec=15

# Outbound only: 127.0.0.1:8001 (RPyC) and api.telegram.org.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/forex-bot/logs /opt/forex-bot/state /opt/forex-bot/data

[Install]
WantedBy=multi-user.target
